# infra/systemd/lac-scraper.timer
[Unit]
Description=Lac — run scraper every 20 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=20min
RandomizedDelaySec=30s
Persistent=true
Unit=lac-scraper.service

[Install]
WantedBy=timers.target
