# infra/systemd/lac-ai.timer
[Unit]
Description=Lac — generate AI commentary once a day at 07:00

[Timer]
OnCalendar=*-*-* 07:00:00
RandomizedDelaySec=5min
Persistent=true
Unit=lac-ai.service

[Install]
WantedBy=timers.target
