# infra/systemd/lac-ai.timer
[Unit]
Description=Lac — tick AI refresher every hour (policy decides if it generates)

[Timer]
# Toutes les heures à xx:55. Le script lit ai_policy et décide lui-même.
OnCalendar=*-*-* *:55:00
RandomizedDelaySec=2min
Persistent=true
Unit=lac-ai.service

[Install]
WantedBy=timers.target
